# rootfs/etc/profile.d/dhyve.sh
#!/bin/sh

echo "DhyveOS version $(cat /etc/version), build $(cat /etc/dhyve)"
docker -v # "Docker version 1.3.0-dev, build ba14ddf-dirty"
