# rootfs/etc/profile.d/dhyve.sh
#!/bin/sh

echo "DhyveOS version $(cat /etc/version)"
docker -v
