# rootfs/var/spool/cron/crontabs/root
*/15 * * * * /usr/bin/sntp -M 128 pool.ntp.org
